# linkedin-prospector.service
[Unit]
Description=LinkedIn Prospector
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/apps/linkedin-prospector
Environment=NODE_ENV=production
Environment=PORT=3462
Environment=HOSTNAME=0.0.0.0
Environment=DATABASE_PATH=/home/ubuntu/apps/linkedin-prospector/data/app.db
Environment=JWT_SECRET=moco-linkedin-prospector-secure-key
Environment=NEXTAUTH_SECRET=moco-linkedin-prospector-secure-key
Environment=NEXTAUTH_URL=https://lp.moco.inc
ExecStart=/usr/bin/npm start
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
